# server-configs/supervisor-archers-london.ini
[program:archers]
command=/srv/http/archers-london/archers/backend/game.py "thebridge" "--location=London, UK" "--maxplayers=8"
directory=/srv/http/archers-london/archers/backend
environment=PATH=/srv/http/archers-london/env/bin
autorestart=true
user=nobody
