# server-configs/supervisor-archers-amsterdam.ini
[program:archers]
command=/srv/http/archers-amsterdam/archers/backend/game.py "thebridge" "--location=Amsterdam, NL" "--maxplayers=10"
directory=/srv/http/archers-amsterdam/archers/backend
environment=PATH=/srv/http/archers-amsterdam/env/bin
autorestart=true
user=nobody
