# server-configs/supervisor-archers-singapore.ini
[program:archers]
command=/srv/http/archers-singapore/archers/backend/game.py "thebridge" "--location=Singapore" "--maxplayers=10"
directory=/srv/http/archers-singapore/archers/backend
environment=PATH=/srv/http/archers-singapore/env/bin
autorestart=true
user=nobody
