# server-configs/supervisor-archers-sf.ini
[program:archers]
command=/srv/http/archers-sf/archers/backend/game.py "thebridge" "--location=San Francisco, USA" "--maxplayers=10"
directory=/srv/http/archers-sf/archers/backend
environment=PATH=/srv/http/archers-sf/env/bin
autorestart=true
user=nobody
